# etc/me.pansen.tunmux.autoconnect.plist
<?xml version="1.0" encoding="UTF-8" ?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>me.pansen.tunmux.autoconnect</string>

    <!-- connect --if-missing is idempotent: it brings the tunnel up if down and
         exits 0 (no-op) when this same profile is already live, so the agent can
         call the binary directly with no wrapper script. -->
    <key>ProgramArguments</key>
    <array>
        <string>@TUNMUX_BIN@</string>
        <string>--debug</string>
        <string>wgconf</string>
        <string>connect</string>
        <string>@CONNECT_FLAG@</string>
        <string>@CONNECT_VALUE@</string>
        <string>--if-missing</string>
    </array>

    <key>RunAtLoad</key>
    <true />
    <key>StartInterval</key>
    <integer>60</integer>   <!-- re-check & reconnect every 60s -->

    <!-- Do NOT set KeepAlive: connect exits 0 after setup, so KeepAlive
         would relaunch it in a tight loop. StartInterval is the right knob. -->

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>   <!-- only in a real GUI login session -->

    <key>StandardOutPath</key>
    <string>@TUNMUX_HOME@/Library/Logs/tunmux-autoconnect.out.log</string>
    <key>StandardErrorPath</key>
    <string>@TUNMUX_HOME@/Library/Logs/tunmux-autoconnect.err.log</string>
</dict>
</plist>
