# etc/me.pansen.tunmux.privileged.plist
<?xml version="1.0" encoding="UTF-8" ?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>me.pansen.tunmux.privileged</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/tunmux</string>
        <string>--debug</string>
        <string>privileged</string>
        <string>--serve</string>
        <string>--authorized-group</string>
        <string>tunmux</string>
    </array>

    <key>RunAtLoad</key>
    <true />
    <key>KeepAlive</key>
    <true />   <!-- genuinely long-running; restart on crash -->

    <key>StandardOutPath</key>
    <string>/var/log/tunmux/privileged.out.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/tunmux/privileged.err.log</string>
</dict>
</plist>
